# pkg/scripts/postinstall
#!/bin/sh
# Dedicated least-privilege accounts the node and reti services run as (neither needs root).
# algod resolves its global config/telemetry path from the account's home, so each
# account needs a real, writable home directory (not /var/empty).
create_hidden_user() {
  user="$1"
  realname="$2"
  home="$3"
  if ! dscl /Local/Default -read "/Users/${user}" >/dev/null 2>&1; then
    NEW_UID=200
    while dscl /Local/Default -list /Users UniqueID | awk '{print $2}' | grep -q "^${NEW_UID}$"; do
      NEW_UID=$((NEW_UID + 1))
    done
    dscl /Local/Default -create "/Users/${user}"
    dscl /Local/Default -create "/Users/${user}" RealName "${realname}"
    dscl /Local/Default -create "/Users/${user}" UserShell /usr/bin/false
    dscl /Local/Default -create "/Users/${user}" UniqueID "$NEW_UID"
    dscl /Local/Default -create "/Users/${user}" PrimaryGroupID 1
    dscl /Local/Default -create "/Users/${user}" IsHidden 1
  fi
  dscl /Local/Default -create "/Users/${user}" NFSHomeDirectory "${home}"
  mkdir -p "${home}"
  chown "${user}" "${home}"
}
create_hidden_user _func-node "FUNC Node" /usr/local/var/func-node
create_hidden_user _func-reti "FUNC Reti" /usr/local/var/func-reti
mkdir -p /usr/local/share/func
FRESH_INSTALL=1
[ -f /Library/LaunchDaemons/func.api.plist ] && FRESH_INSTALL=0
cp /opt/func/Templates/func.api.plist /Library/LaunchDaemons
launchctl bootstrap system /Library/LaunchDaemons/func.api.plist
sleep 2
if [ "$FRESH_INSTALL" -eq 1 ]; then
  console_uid=$(stat -f %u /dev/console)
  if [ "$console_uid" -gt 0 ]; then
    launchctl asuser "$console_uid" sudo -u "#$console_uid" open http://localhost:3536
  fi
fi
exit 0