# pkg/scripts/postinstall
#!/bin/sh
mkdir -p /usr/local/share/func
cp /opt/func/Templates/func.api.plist /Library/LaunchDaemons
launchctl bootstrap system /Library/LaunchDaemons/func.api.plist
sleep 2
open http://localhost:3536
exit 0